# home-manager/domain/autostart/org.kde.polybar.start.desktop
[Desktop Entry]
Name=polybar launch
Exec=fish -c polybar-launch
Icon=system-software-update
Type=Application
NoDisplay=true
OnlyShowIn=KDE
